# rootfs/etc/services.d/gmc320/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start gmc320 service
# ==============================================================================

PORT="$(bashio::config 'port')"
BAUDRATE="$(bashio::config 'baudrate')"
DELAY="$(bashio::config 'repeat')"

set +e

bashio::log.info "Starting gmc3xx, using port ${PORT} and baudrate ${BAUDRATE}"

# Get the device serial number
RESPONSE=$(gmc320 ${PORT} ${BAUDRATE})
bashio::log.info "response: ${RESPONSE}"
SERIAL=$(echo "${RESPONSE}" | jq -r .attributes.serial)
bashio::log.info "Device serial number: ${SERIAL}"


# Send the readings to http://supervisor/core/api
bashio::log.info "Start sending data"
while true; do
  bashio::log.info "Requesting data, using port ${PORT} and baudrate ${BAUDRATE}"
	DATA=$(gmc320 ${PORT} ${BAUDRATE})
	bashio::log.info "response: ${DATA}"
	#mosquitto_pub -u ${MQTT_USER} -P ${MQTT_PASSWORD} -p ${MQTT_PORT} -h ${MQTT_SERVER} -t "homeassistant/sensor/gmc3xx_${SERIAL}" -m "${DATA}"
	HTTP_CODE=$(curl -s -X POST -o response.txt -w "%{http_code}" \
    -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
    -H "Content-Type: application/json" \
    -d "${DATA}" \
    "http://supervisor/core/api/states/sensor.gmc3xx_${SERIAL}")
  HTTP_RESPONSE=$(cat response.txt)
  bashio::log.info "HA API HTTP code: ${HTTP_CODE}"
  bashio::log.info "HA API HTTP Response: ${HTTP_RESPONSE}"
	sleep ${DELAY}
done
